# .cache/sigma/T1036_004.json
{
  "technique_id": "T1036.004",
  "cached_at": "2026-04-01T05:35:58.843635+00:00",
  "rules": [
    {
      "title": "Scheduled Task Creation Masquerading as System Processes",
      "status": "experimental",
      "level": "high",
      "logsource": "process_creation / windows",
      "condition_summary": "all of selection_*",
      "description": "Detects the creation of scheduled tasks that involve system processes, which may indicate malicious actors masquerading as or abusing these processes to execute payloads or maintain persistence.",
      "tags": [
        "attack.privilege-escalation",
        "attack.execution",
        "attack.persistence",
        "attack.t1053.005",
        "attack.defense-evasion",
        "attack.t1036.004",
        "attack.t1036.005"
      ],
      "url": "https://github.com/SigmaHQ/sigma/blob/2f84ca2f1652977cd59f48b6556dfc3f6f14fb5a/rules/windows/process_creation/proc_creation_win_schtasks_system_process.yml"
    }
  ]
}